# packaging/kakku-niri-settings/kakku-niri-settings.install
_enable_defaults() {
  systemctl enable power-profiles-daemon.service 2>/dev/null || true
  systemctl --global enable dms.service 2>/dev/null || true
  systemctl --global enable dsearch.service 2>/dev/null || true

  # Override os-release with KakkuOS branding
  if [ -f /usr/share/kakku/os-release ]; then
    cp /usr/share/kakku/os-release /usr/lib/os-release
  fi

  if command -v kakku-limine-defaults >/dev/null 2>&1; then
    kakku-limine-defaults || true
  fi

  # Remove CachyOS welcome app autostart
  rm -f /etc/xdg/autostart/cachyos-hello.desktop
}

post_install() {
  _enable_defaults
  _apply_missing_configs
}

post_upgrade() {
  _enable_defaults
  _apply_missing_configs
}

_apply_missing_configs() {
  local skel="/etc/skel"

  for user_home in /home/*; do
    [ -d "$user_home" ] || continue
    local username
    username=$(basename "$user_home")
    id "$username" &>/dev/null || continue

    # Seed missing files inside .config directories without replacing user edits.
    for dir in DankMaterialShell fastfetch git lazygit menus mpv starship yazi zathura fish niri niri-screensaver ghostty nvim kakku github-copilot claude; do
      if [ -d "$skel/.config/$dir" ]; then
        _copy_missing_tree_files "$skel/.config/$dir" "$user_home/.config/$dir" "$username"
      fi
    done

    # Shell rc files are individual Kakku-owned defaults; seed only if missing.
    if [ -f "$skel/.bashrc" ]; then
      _copy_missing_file "$skel/.bashrc" "$user_home/.bashrc" "$username"
    fi

    if [ -f "$skel/.zshrc" ]; then
      _copy_missing_file "$skel/.zshrc" "$user_home/.zshrc" "$username"
    fi

    if command -v kakku-browser-theme >/dev/null 2>&1 && command -v runuser >/dev/null 2>&1; then
      runuser -u "$username" -- env HOME="$user_home" XDG_CONFIG_HOME="$user_home/.config" XDG_CACHE_HOME="$user_home/.cache" kakku-browser-theme >/dev/null 2>&1 || true
    fi

    if command -v kakku-vscode-theme >/dev/null 2>&1 && command -v runuser >/dev/null 2>&1; then
      runuser -u "$username" -- env HOME="$user_home" XDG_CONFIG_HOME="$user_home/.config" kakku-vscode-theme >/dev/null 2>&1 || true
    fi
  done
}

_copy_missing_file() {
  local source="$1"
  local target="$2"
  local username="$3"

  if [ -e "$target" ]; then
    return 0
  fi

  mkdir -p "$(dirname "$target")"
  cp -a "$source" "$target"
  chown "$username:$username" "$target"
}

_copy_missing_tree_files() {
  local source="$1"
  local target="$2"
  local username="$3"
  local source_path
  local relative_path
  local target_path

  if [ ! -d "$target" ]; then
    mkdir -p "$target"
    chown "$username:$username" "$target"
  fi

  find "$source" -type d | while IFS= read -r source_path; do
    relative_path="${source_path#$source}"
    target_path="$target$relative_path"

    if [ -d "$target_path" ]; then
      continue
    fi

    mkdir -p "$target_path"
    chown "$username:$username" "$target_path"
  done

  find "$source" -type f | while IFS= read -r source_path; do
    relative_path="${source_path#$source/}"
    target_path="$target/$relative_path"

    mkdir -p "$(dirname "$target_path")"
    if [ -e "$target_path" ]; then
      continue
    fi

    cp -a "$source_path" "$target_path"
    chown "$username:$username" "$target_path"
  done
}
