# system/systemd/user/kakku-idle.service
[Unit]
Description=KakkuOS idle screensaver and DMS lock handoff
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/kakku-idle
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
