# packaging/kakku-hyprland-settings/kakku-hyprland-settings.install
post_install() {
  # Override os-release with KakkuOS branding
  if [ -f /usr/share/kakku/os-release ]; then
    cp /usr/share/kakku/os-release /usr/lib/os-release
  fi

  # Enable greetd as login manager (disable others)
  systemctl disable sddm.service 2>/dev/null || true
  systemctl disable ly.service 2>/dev/null || true
  systemctl enable greetd.service 2>/dev/null || true

  # Apply configs to existing users
  _apply_configs
}

post_upgrade() {
  post_install
}

_apply_configs() {
  local skel="/etc/skel"

  for user_home in /home/*; do
    [ -d "$user_home" ] || continue
    local username=$(basename "$user_home")
    id "$username" &>/dev/null || continue

    # Copy .config directories
    for dir in hypr waybar rofi alacritty fastfetch mako lazygit yazi btop; do
      if [ -d "$skel/.config/$dir" ]; then
        mkdir -p "$user_home/.config/$dir"
        cp -r "$skel/.config/$dir/." "$user_home/.config/$dir/"
        chown -R "$username:$username" "$user_home/.config/$dir"
      fi
    done

    # Copy .zshrc
    if [ -f "$skel/.zshrc" ]; then
      cp "$skel/.zshrc" "$user_home/.zshrc"
      chown "$username:$username" "$user_home/.zshrc"
    fi
  done
}
